# python/catbus/catbus_directory.service
[Unit]
Description=Catbus Directory Server

[Service]
Type=simple
ExecStart=/usr/local/bin/catbus_directory

[Install]
WantedBy=multi-user.target