# python/catbus/catbus_directory.conf

[program:catbus_directory]
command=catbus_directory
directory=/var/log/chromatron/
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=TERM
user=catbus
